# roles/backup-duplicity-docker/templates/etc/systemd/system/duplicity-docker-tpl.timer.j2
[Unit]
Description=Duplicity: {{ item }}

[Timer]
OnCalendar={{ duplicity[item] }}
Persistent=True

[Install]
WantedBy=timers.target
