# roles/bluetooth/files/etc/udev/rules.d/80-bt-headset.rules
# Managed by Ansible
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/usr/local/sbin/a2dp-fix-wrapper"
